# scripts/service/container.target
[Unit]
Description=Container target
AllowIsolate=yes
