# packer/rootfs/etc/rc.local
#!/bin/sh

IP_ADDRESS=$(hostname -I ||:)
if [ -n "${IP_ADDRESS?}" ]; then
	printf '%s\n' "My IP address is: ${IP_ADDRESS:?}"
fi

if [ -e /opt/vc/bin/tvservice ] && ! /opt/vc/bin/tvservice -s | grep -Eq 'HDMI|DVI'; then
	/opt/vc/bin/tvservice -o
fi

exit 0
