# packer/rootfs/etc/rc.local
#!/bin/sh

set -eu

IP_ADDRESS=$(hostname -I ||:)
if [ -n "${IP_ADDRESS?}" ]; then
	printf '%s\n' "My IP address is: ${IP_ADDRESS:?}"
fi

if [ -c /dev/vchiq ] && ! tvservice --status | grep -Eq 'HDMI|DVI'; then
	tvservice --off
fi

exit 0
